# etc/supervisord.conf
[supervisord]
logfile=%(ENV_PH_HOME_DIR)s/var/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=%(ENV_PH_HOME_DIR)s/var/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[inet_http_server]
port=127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001


[group:apps]
programs = polls_server, polls_client


[program:polls_server]
directory = %(ENV_PH_HOME_DIR)s/polls_server
command = %(ENV_PH_HOME_DIR)s/bin/gunicorn -c %(ENV_PH_HOME_DIR)s/etc/gunicorn_server.conf polls_server.wsgi:application
stdout_logfile = %(ENV_PH_HOME_DIR)s/var/log/server-stdout.log
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8,HOME="/home/alex/workspace/projects/python/polls_hub/polls_server",DJANGO_SETTINGS_MODULE="polls_server.settings"
autostart=true


[program:polls_client]
directory = %(ENV_PH_HOME_DIR)s/polls_client
command = %(ENV_PH_HOME_DIR)s/bin/gunicorn -c %(ENV_PH_HOME_DIR)s/etc/gunicorn_client.conf polls_client.wsgi:application
redirect_stderr = true
stdout_logfile = %(ENV_PH_HOME_DIR)s/var/log/client-stdout.log
autostart = true
